# vsftpd.target
[Unit]
Description=FTP daemon
Documentation=man:vsftpd
